# deploy/systemd/olx-dashboard.service
[Unit]
Description=OLX Deals dashboard (local web UI, reachable over Tailscale)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/daniel/Projects/OLX
ExecStart=/usr/bin/python3 -m olxdeals.dashboard --db /home/daniel/Projects/OLX/olxdeals.db --config /home/daniel/Projects/OLX/searches.yaml --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
